# .config/systemd/user/blueberry.service
[Unit]
Description=Bluetooth Tray Daemon

[Service]
Type=forking
ExecStart=/usr/bin/blueberry-tray

[Install]
WantedBy=default.target
