# .config/systemd/user/blueberry.service
[Unit]
Description=Bluetooth Tray Daemon

[Service]
ExecStart=/usr/bin/blueberry

[Install]
WantedBy=default.target
